# deploy/dismissal-health.service
[Unit]
Description=Dismissal Health Endpoint — lightweight HTTP status server
Documentation=https://github.com/joker5914/Dismissal
# Don't gate on network-online.target — see dismissal-scanner.service.
# The health endpoint binds to 0.0.0.0 and is reachable as soon as
# NetworkManager activates an interface; waiting for "fully online"
# just delays the LAN diagnostic.
After=NetworkManager.service dismissal-setup-portal.service
Wants=NetworkManager.service
# Wait for the captive portal to finish provisioning WiFi before exposing
# the device's health endpoint.
ConditionPathExists=/var/lib/dismissal/.wifi-provisioned

[Service]
User=dismissal
Group=dismissal
WorkingDirectory=/opt/dismissal/Backend
# Optional — scanner config comes from scanner_config.py in the repo.
EnvironmentFile=-/opt/dismissal/Backend/.env
ExecStart=/opt/dismissal/venv/bin/python /opt/dismissal/Backend/dismissal_health.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dismissal-health
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
