# deploy/dismissal-health.service
[Unit]
Description=Dismissal Health Endpoint — lightweight HTTP status server
Documentation=https://github.com/joker5914/Dismissal
After=network-online.target dismissal-setup-portal.service
Wants=network-online.target
# Wait for the captive portal to finish provisioning WiFi before exposing
# the device's health endpoint.
ConditionPathExists=/var/lib/dismissal/.wifi-provisioned

[Service]
User=dismissal
Group=dismissal
WorkingDirectory=/opt/dismissal/Backend
# Optional — scanner config comes from scanner_config.py in the repo.
EnvironmentFile=-/opt/dismissal/Backend/.env
ExecStart=/opt/dismissal/venv/bin/python /opt/dismissal/Backend/dismissal_health.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dismissal-health
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
