# deploy/dismissal-scanner.service
[Unit]
Description=Dismissal Licence Plate Scanner
Documentation=https://github.com/joker5914/Dismissal
# Wait for the network to be actually online before starting.  Without
# this, the scanner races WiFi association and its first Firebase call
# fails with "Temporary failure in name resolution," then back-off
# kicks in and the device doesn't reach the cloud for 1–2 minutes.
#
# The 2-minute boot we saw before came from systemd-networkd-wait-
# online sitting on its full timeout — that service is irrelevant on a
# NetworkManager-managed Pi.  install.sh + update.sh now mask it, so
# network-online.target is satisfied as soon as NetworkManager-wait-
# online completes (typically ~6 s after NM activates the WiFi link).
After=network-online.target
Wants=network-online.target
After=dismissal-watchdog.service dismissal-setup-portal.service
# An unprovisioned device puts up the captive portal AP instead of the
# scanner.  The portal writes this marker when WiFi is set; until then
# the scanner unit silently stays inactive (ConditionPathExists treats a
# failed condition as "condition not met", not as a failure — no log
# spam, no restart loop).  Once the portal succeeds it explicitly runs
# `systemctl start dismissal-scanner`, so we don't need a reboot.
ConditionPathExists=/var/lib/dismissal/.wifi-provisioned
# StartLimit belongs in [Unit], not [Service] — silences the
# "Unknown key 'StartLimitIntervalSec'" warning we saw on Trixie.
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
# ---------------------------------------------------------------------------
# Type=notify — required for sd_notify(READY=1), STOPPING=1, and
# WATCHDOG=1 keepalives to work.  NOTIFY_SOCKET and WATCHDOG_USEC are only
# exported by systemd when this type is set.
# ---------------------------------------------------------------------------
Type=notify
User=dismissal
Group=dismissal
# WorkingDirectory points at the OTA-swappable release.  Whatever
# /opt/dismissal/current/ resolves to is the active release.  The
# symlink is flipped atomically by deploy/firmware_swap.sh; on a
# fresh install the migration in update.sh sets it up so this path
# always resolves to a real Backend/ directory.
WorkingDirectory=/opt/dismissal/current/Backend
# ---------------------------------------------------------------------------
# Thermal-conscious defaults for the Pi 5.  Each Environment= value is
# overridden by the same key in /opt/dismissal/Backend/.env when set,
# so operators can dial back up for diagnostics without editing this unit.
#
#   SCANNER_RESOLUTION=1280x720  — 56% of 1080p pixels through every
#                                  motion-gate, ML, and JPEG-encode path.
#   SCANNER_CAMERA_FPS=10        — match the pipeline FPS_CAP so the
#                                  libcamera ISP isn't producing 30 fps
#                                  the scanner immediately throws away.
#   SCANNER_DEBUG_PORT=0         — disables the LAN MJPEG view; the
#                                  encoder thread copies+annotates+JPEGs
#                                  every motion-positive frame even with
#                                  zero clients connected.  Set back to
#                                  8081 in .env when actively debugging.
#   OMP/OPENBLAS/MKL_NUM_THREADS=2 — leave a core idle for cooling
#                                  headroom; the Pi 5 throttles on
#                                  per-core temps, not just package.
# ---------------------------------------------------------------------------
Environment=SCANNER_RESOLUTION=1280x720
Environment=SCANNER_CAMERA_FPS=10
Environment=SCANNER_DEBUG_PORT=0
Environment=OMP_NUM_THREADS=2
Environment=OPENBLAS_NUM_THREADS=2
Environment=MKL_NUM_THREADS=2
# HailoRT writes its persistent state to ~/.hailo on first init.  The
# dismissal user's home is /opt/dismissal, but ProtectSystem=strict +
# the ReadWritePaths allowlist below makes that read-only — Hailo init
# fails (CHECK failed - Failed to create directory /opt/dismissal/.hailo)
# and the scanner silently falls back from plate_hailo to plate_yolo on
# CPU.  Pointing HOME at the already-writable StateDirectory lets Hailo
# create /var/lib/dismissal/.hailo/ on its own, no other config changes.
Environment=HOME=/var/lib/dismissal

# Leading `-` makes the EnvironmentFile optional — service starts fine with
# no .env at all, since scanner_config.py provides the required defaults.
# .env is used to override tuning values (FPS, plate length, etc.) and
# takes precedence over the Environment= lines above.
EnvironmentFile=-/opt/dismissal/Backend/.env

# ---------------------------------------------------------------------------
# State directory — systemd creates /var/lib/dismissal owned by
# dismissal:dismissal before the service starts.  Used for the SQLite outbox.
# ---------------------------------------------------------------------------
StateDirectory=dismissal
StateDirectoryMode=0750

ExecStart=/opt/dismissal/venv/bin/python /opt/dismissal/current/Backend/dismissal.py
ExecReload=/bin/kill -HUP $MAINPID

# ---------------------------------------------------------------------------
# Restart policy
# ---------------------------------------------------------------------------
Restart=on-failure
RestartSec=5

# ---------------------------------------------------------------------------
# Logging
# ---------------------------------------------------------------------------
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dismissal-scanner

# ---------------------------------------------------------------------------
# Hardening
# ProtectSystem=strict makes the entire filesystem read-only except for
# paths listed in ReadWritePaths.  PrivateDevices=false and the DeviceAllow
# lines restore access to the camera and Coral USB accelerator.
# ---------------------------------------------------------------------------
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

ReadWritePaths=/opt/dismissal/Backend/debug_frames \
               /opt/dismissal/Backend/logs \
               /opt/dismissal/current/Backend/debug_frames \
               /opt/dismissal/current/Backend/logs \
               /opt/dismissal/releases \
               /var/lib/dismissal \
               /tmp

# Camera and Coral USB access — PrivateDevices=false is required so /dev
# is not replaced with a minimal tmpfs that hides hardware devices.
#
# Pi 5 / libcamera exposes far more device nodes than Pi 4 did:
#   /dev/video* and /dev/v4l-subdev*  (sensors + ISP chains)  → char-video4linux
#   /dev/media0..N                    (media controller pipes) → char-media
#   /dev/dma_heap/system              (dma-buf for zero-copy)  → char-dma_heap
# Listing each node by path missed the new ones and caused libcamera to log
# "Operation not permitted" errors on /dev/media2, /dev/video19, and friends,
# so Picamera2 init failed and the scanner fell back to OpenCV (which also
# couldn't open the camera). Allow the whole char classes instead.
PrivateDevices=false
DeviceAllow=char-video4linux rw
DeviceAllow=char-media rw
DeviceAllow=char-dma_heap rw
# Coral USB Accelerator — usb_device char class covers /dev/bus/usb/*/*
DeviceAllow=char-usb_device rw
# Hailo AI HAT+ — the kernel module exposes /dev/hailo0 with a custom
# character device major; there's no standard "char-hailo" class so we
# allow the path explicitly.  Without this the dismissal user under
# systemd's device cgroup can't open the chip even though /dev/hailo0
# is world-writable, and PlateHailoDetector init bails with
# HAILO_DRIVER_OPERATION_FAILED(36) → silent fallback to ONNX/CPU.
DeviceAllow=/dev/hailo0 rw

# ---------------------------------------------------------------------------
# systemd watchdog — scanner pings WATCHDOG=1 every WatchdogSec/2 seconds.
# If it goes silent for WatchdogSec, systemd kills and restarts the service.
# ---------------------------------------------------------------------------
WatchdogSec=60
NotifyAccess=main

# ---------------------------------------------------------------------------
# Resource limits — sized for Pi 5 (4–8 GB RAM, 4-core Cortex-A76)
# CPUQuota of 300% allows up to 3 of 4 cores at full speed.
# MemoryMax of 1536M leaves headroom for the OS, watchdog, and health services.
# ---------------------------------------------------------------------------
CPUQuota=300%
MemoryMax=1536M
TasksMax=128

[Install]
WantedBy=multi-user.target
